# .docker/supervisor/laravel-workers.ini
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/site/artisan queue:work --timeout=600 --tries=3 --max-jobs=1000
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/laravel-queue.log
stopwaitsecs=3600
stdout_logfile_maxbytes=5MB
