# .docker/supervisor/unit.conf
[supervisord]
logfile = /tmp/supervisord.log
pidfile = /tmp/supervisord.pid
nodaemon=true
user=root

[program:unit]
command=unitd --no-daemon --control unix:/var/run/control.unit.sock
pidfile=/run/unit.pid
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0